# crontab/reconcile_crontab
# how to insert above lines to crontab
# $ crontab -e

# The first asterisk is for specifying the minute of the run (0-59)
# The second asterisk is for specifying the hour of the run (0-23)
# The third asterisk is for specifying the day of the month for the run (1-31)
# The fourth asterisk is for specifying the month of the run (1-12)
# The fifth asterisk is for specifying the day of the week (where Sunday is equal to 0, up to Saturday is equal to 6)

# please copy following lines to crontab

LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/local/lib:/usr/local/pgsql/lib

0   1   *   *   *   /Users/winnuayi/Projects/dev/backend/reconcile-aj/run.sh
0   12   *   *   *   /Users/winnuayi/Projects/dev/backend/reconcile-aj/run.sh
